# serve/glm-planner.service
[Unit]
Description=llama-server: GLM-4.7-Flash planner (thinking on) OpenAI API on the LAN
After=network-online.target
Wants=network-online.target
# only one model fits in 24GB — starting the planner stops the coder and vice-versa
Conflicts=qwen3-coder.service

[Service]
Type=simple
User=psc
Group=psc
WorkingDirectory=/home/psc/dev/llm
Environment=HOME=/home/psc
Environment=PATH=/home/psc/miniconda3/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=LD_LIBRARY_PATH=/home/psc/miniconda3/lib
ExecStart=/usr/bin/bash /home/psc/dev/llm/serve/serve.sh glm-planner
Restart=on-failure
RestartSec=5

# NOTE: install but do NOT `enable` this one — it's started on demand (the coder is the default
# that auto-starts at boot).
[Install]
WantedBy=multi-user.target
