# serve/llama-swap.service
[Unit]
Description=llama-swap: one OpenAI endpoint on the LAN, auto-swaps llama.cpp models on request
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=psc
Group=psc
WorkingDirectory=/home/psc/dev/llm
Environment=HOME=/home/psc
Environment=PATH=/home/psc/miniconda3/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
# llama-server's cudart rpath is baked in, but keep conda libs reachable for any upstream.
Environment=LD_LIBRARY_PATH=/home/psc/miniconda3/lib
# --watch-config: editing/regenerating serve/llama-swap.yaml hot-reloads without a restart.
ExecStart=/home/psc/opt/llama-swap/llama-swap --config /home/psc/dev/llm/serve/llama-swap.yaml --listen 0.0.0.0:8080 --watch-config
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
