# serve/qwen3-coder.service
[Unit]
Description=llama-server: qwen3-coder OpenAI-compatible API on the LAN
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=psc
Group=psc
WorkingDirectory=/home/psc/dev/llm
# llama-server has the conda cudart rpath baked in; these just ensure serve.sh's python3
# (to read models.toml) and libs resolve under systemd's minimal environment.
Environment=HOME=/home/psc
Environment=PATH=/home/psc/miniconda3/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=LD_LIBRARY_PATH=/home/psc/miniconda3/lib
ExecStart=/usr/bin/bash /home/psc/dev/llm/serve/serve.sh qwen3-coder
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
